# deploy/ec2/opensandbox-gitserver.service
[Unit]
Description=OpenSandbox Git Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/opensandbox-gitserver
Restart=always
RestartSec=5
Environment=HOME=/root
Environment=OPENSANDBOX_GIT_PORT=3000
Environment=OPENSANDBOX_GIT_REPO_ROOT=/data/opensandbox-repos
Environment=OPENSANDBOX_GIT_DOMAIN=git.opensandbox.ai
Environment=OPENSANDBOX_SECRETS_ARN=arn:aws:secretsmanager:us-east-2:739940681129:secret:opensandbox/server-vtN2Ez
EnvironmentFile=-/etc/opensandbox/server.env
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
